# systemSettings/etc/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f -- \\u' --noclear --autologin lukaszpe %I $TERM
Type=simple
#Type=oneshot
Environment=XDG_SESSION_TYPE=wayland
